# deploy/alloy/alloy.systemd.service
[Unit]
Description=Grafana Alloy — MoE Sovereign observability collector
Documentation=https://grafana.com/docs/alloy/
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=alloy
Group=alloy
EnvironmentFile=-/etc/default/alloy
ExecStart=/usr/local/bin/alloy run /etc/alloy/config.river \
          --server.http.listen-addr=127.0.0.1:12345 \
          --storage.path=/var/lib/alloy
Restart=on-failure
RestartSec=10

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/alloy
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Needs journal read access for loki.source.journal
SupplementaryGroups=systemd-journal

[Install]
WantedBy=multi-user.target
